# templates/openrc/autodoc-github-research-scheduler
#!/sbin/openrc-run

name="Autodoc GitHub research Scheduler"
description="Scheduler canonique durable pour les recherches GitHub Autodoc"

: "${AUTODOC_ROOT:=/home/eric/projet/git/autodoc}"
: "${AUTODOC_PYTHON:=/home/eric/python/bin/python}"
: "${AUTODOC_RUNTIME_CONFIG:=${AUTODOC_ROOT}/.var/config/love_installed_runtime.ini}"
: "${AUTODOC_SERVICE_USER:=eric}"
: "${AUTODOC_GITHUB_RESEARCH_OPENRC_SERVICE_FACTORY:=}"

command="${AUTODOC_PYTHON}"
command_args="${AUTODOC_ROOT}/tools/run_github_research_love_openrc_scheduler_0287.py --factory ${AUTODOC_GITHUB_RESEARCH_OPENRC_SERVICE_FACTORY} --config ${AUTODOC_RUNTIME_CONFIG} --format summary"
command_user="${AUTODOC_SERVICE_USER}"
directory="${AUTODOC_ROOT}"
command_background="yes"
pidfile="/run/autodoc/github-research-scheduler.pid"
start_stop_daemon_args="--make-pidfile"

start_pre() {
    if [ -z "${AUTODOC_GITHUB_RESEARCH_OPENRC_SERVICE_FACTORY}" ]; then
        eerror "AUTODOC_GITHUB_RESEARCH_OPENRC_SERVICE_FACTORY est absent"
        return 1
    fi
    if [ ! -x "${AUTODOC_PYTHON}" ]; then
        eerror "Python Autodoc introuvable: ${AUTODOC_PYTHON}"
        return 1
    fi
    if [ ! -f "${AUTODOC_RUNTIME_CONFIG}" ]; then
        eerror "Configuration runtime introuvable: ${AUTODOC_RUNTIME_CONFIG}"
        return 1
    fi
    checkpath --directory --mode 0750 --owner "${AUTODOC_SERVICE_USER}" /run/autodoc
}

depend() {
    need localmount postgresql-18 qdrant
    after net
}
